# crontab
*/15 * * * * root /var/script/sync.sh >> /var/log/cron.log 2>&1
